# etc/systemd/balerd.service.in
# README
# ======
# Please copy this file to /etc/systemd/system/ before
# applying any of the following commands.
#
# To make balerd starts automatically at system startup time:
#   systemctl enable balerd.service
#
# To disable automatic startup:
#   systemctl disable balerd.service
#
# To manually start the service:
#   systemctl start balerd.service
#
# To check status:
#   systemctl status balerd.service
#
# To stop the service:
#   sytemctl stop balerd.service

[Unit]
Description = Baler Daemon
Documentation = http://ovis.ca.sandia.gov

[Service]
Type = forking
EnvironmentFile = @sysconfdir@/baler/balerd.env
Environment = "HOSTNAME=$(hostname)"
ExecStart = @bindir@/balerd \
		-s ${BALERD_STORE} \
		-S ${BALERD_STORE_PLUGIN} \
		-C ${BALERD_PLUGIN_CONFIG} \
		-v ${BALERD_LOG_LEVEL} \
		-I ${BALERD_IN_Q_WORKERS} \
		-O ${BALERD_OUT_Q_WORKERS}

[Install]
WantedBy = default.target
